# .service
[Unit]
Description=Gunicorn instance to serve flask server
After=network.target

StartLimitIntervalSec=500
StartLimitBurst=5

[Service]
User=max
Group=max
WorkingDirectory=/srv/git/server
Environment="PATH=/srv/git/server/env/bin"
ExecStart=/srv/git/server/env/bin/gunicorn -w 3 -k 'gevent' --bind 0.0.0.0:5000 'app:create_app()'
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target